# udev/rules/99-kensington-orbit.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="047d", ATTRS{idProduct}=="2048", RUN+="/opt/udev/kensington-trackball"
